# dist/ping_exporter.service
[Unit]
Description=Ping Exporter
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/ping_exporter --config.path=/etc/ping_exporter/config.yml
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target
