# dist/ping_exporter.service
[Unit]
Description=Ping Exporter
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/ping_exporter -web.listen-address 127.0.0.1:9427 8.8.8.8 9.9.9.9
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target
